# docs/systemd-example/pajbot2-discord.service.d/10-bot-token.conf
[Service]
Environment=PAJBOT2_DISCORD_BOT_TOKEN=
