# docs/systemd-example/pajbot2-discord.service.d/10-postgres.conf
[Service]
Environment=PAJBOT2_DISCORD_BOT_SQL_DSN="host=/var/run/postgresql database=pajbot2-discord"
