# docs/systemd-example/pajbot2-discord.service.d/10-slash-commands.conf
[Service]
# Enable all slash commands in guilds 12345 and 67890
Environment=PAJBOT2_DISCORD_BOT_SLASH_COMMAND_GUILD_IDS="12345,67890"
